# src/github.com/TheRealBill/redskull/docker/supervisord/redskull.conf
[program:redskull]
command=/usr/redskull/redskull 
directory=/usr/redskull
