# docker/supervisord/consul.conf
[program:consul]
command=/usr/bin/consul agent -join 172.17.42.1 -data-dir=/tmp/consul -config-dir=/etc/consul.d
