# docker/supervisord/sentinel.conf
[program:sentinel]
command=/usr/sbin/redis-server /etc/redis/sentinel.conf --sentinel
